# agent/launchd/com.pi-matrix.hermes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pi-matrix.hermes</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/hermes</string>
    <string>--config</string>
    <string>/etc/pi-matrix/hermes-config.yaml</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/var/log/pi-matrix/hermes.log</string>

  <key>StandardErrorPath</key>
  <string>/var/log/pi-matrix/hermes.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
